# .config/autostart/pause-on-lock.desktop
[Desktop Entry]
Type=Application
Exec=pause-on-lock.sh --playerctl
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Pause On lock
